# cron.d/mysb_dbmf3
# replace [*] by the good values

00 06 * * * www-data cd [path_to_phpmysandbox] && /usr/bin/php script.php mod=dbmf3 tpl=notifications spw=[script_password]
